# systemd/waydroid-container.service
[Unit]
Description=Waydroid Container

[Service]
BusName=id.waydro.Container
ExecStart=/usr/bin/waydroid -w container start

[Install]
WantedBy=multi-user.target
