# systemd/waydroid-session.service
[Unit]
Description=Waydroid Session

[Service]
ExecStart=/usr/bin/waydroid session start

[Install]
WantedBy=multi-user.target
