# service_managers/s6/waydroid-srv/run
#!/bin/execlineb -P
fdmove -c 2 1 

foreground {
	nft -f /etc/waydroid/lxc.conf -f /etc/waydroid/nftables.rules
}

exec /usr/bin/waydroid -w container start

